# support/systemd/nextepc-pcrfd.service.in
[Unit]
Description=NextEPC PCRF Daemon
BindTo=mongodb.service
After=networking.service mongodb.service

[Service]
Type=simple

User=nextepc
Group=nextepc

RuntimeDirectory=nextepc-pcrfd
PIDFile=@LIB_DIR@/run/nextepc-pcrfd/pid

Restart=always
ExecStart=@BIN_DIR@/nextepc-pcrfd -f @SYSCONF_DIR@/nextepc/pcrf.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
