# support/systemd/nextepc-mmed.service.in
[Unit]
Description=NextEPC MME Daemon
After=networking.service

[Service]
Type=simple

User=nextepc
Group=nextepc

RuntimeDirectory=nextepc-mmed
PIDFile=@LIB_DIR@/run/nextepc-mmed/pid

Restart=always
ExecStart=@BIN_DIR@/nextepc-mmed -f @SYSCONF_DIR@/nextepc/mme.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
