# support/systemd/nextepc-hssd.service.in
[Unit]
Description=NextEPC HSS Daemon
After=networking.service

[Service]
Type=simple
Restart=always
ExecStart=@BIN_DIR@/nextepc-hssd -f @SYSCONF_DIR@/nextepc/hss.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
