# support/systemd/nextepc-hssd.service.in
[Unit]
Description=NextEPC HSS Daemon
After=networking.service mongodb.service

[Service]
Type=simple

User=nextepc
Group=nextepc

RuntimeDirectory=nextepc-hssd
PIDFile=@LIB_DIR@/run/nextepc-hssd/pid

Restart=always
ExecStartPre=/bin/sleep 5
ExecStart=@BIN_DIR@/nextepc-hssd -f @SYSCONF_DIR@/nextepc/hss.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
