# support/systemd/nextepc-pgwd.service.in
[Unit]
Description=NextEPC S-GW Daemon
After=networking.service

[Service]
Type=simple
Restart=always
ExecStart=@BIN_DIR@/nextepc-pgwd -f @SYSCONF_DIR@/nextepc/pgw.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
