# support/systemd/nextepc-pcrfd.service.in
[Unit]
Description=NextEPC PCRF Daemon
After=networking.service mongodb.service

[Service]
Type=simple
Restart=always
ExecStartPre=/bin/sleep 2
ExecStart=@BIN_DIR@/nextepc-pcrfd -f @SYSCONF_DIR@/nextepc/pcrf.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
