# configs/systemd/open5gs-smsfd.service.in
[Unit]
Description=Open5GS SMSF Daemon
After=network-online.target

[Service]
Type=simple

User=open5gs
Group=open5gs

Restart=always
ExecStart=@bindir@/open5gs-smsfd -c @sysconfdir@/open5gs/smsf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
