# configs/systemd/open5gs-nrfd.service.in
[Unit]
Description=Open5GS NRF Daemon
Wants=mongodb.service mongod.service

[Service]
Type=simple

User=open5gs
Group=open5gs

Restart=always
ExecStart=@bindir@/open5gs-nrfd -c @sysconfdir@/open5gs/nrf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
