# configs/systemd/open5gs-pcfd.service.in
[Unit]
Description=Open5GS PCF Daemon
After=networking.service

[Service]
Type=simple

User=open5gs
Group=open5gs

Restart=always
ExecStart=@bindir@/open5gs-pcfd -c @sysconfdir@/open5gs/pcf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
