# configs/systemd/open5gs-upfd.service.in
[Unit]
Description=Open5GS UPF Daemon
After=network-online.target systemd-networkd.service
Wants=systemd-networkd.service

[Service]
Type=simple

User=open5gs
Group=open5gs

Restart=always
ExecStart=@bindir@/open5gs-upfd -c @sysconfdir@/open5gs/upf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
